# ansible/playbooks/install_desktop/roles/backlight/files/80-power.rules
SUBSYSTEM=="power_supply", ACTION=="change", ENV{POWER_SUPPLY_NAME}=="AC", RUN+="/usr/bin/xbacklight 30"
